# config/crontab
# ----------
# Scheduled jobs for lecture capture system
# ----------
#M H D M W Command

# Start the system at reboot
@reboot /home/capture/Camera/start

# Upload any completed videos
08,38 *  *  *  * /home/capture/Camera/upload

 
